# cookbooks/sickbeard/templates/default/sickbeard.service.erb
[Unit]
Description=Sickbeard indexer
After=network.target

[Service]
User=<%= node['sickbeard']['user'] %>
Group=<%= node['sickbeard']['user'] %>
WorkingDirectory=/usr/local/sickbeard
ExecStart=/usr/bin/python SickBeard.py --config=/usr/local/sickbeard/config.ini

[Install]
WantedBy=multi-user.target
